# scan-monitor.service
[Unit]
Description=Chain monitor

[Service]
ExecStart=/root/scan-monitor/start.sh
Restart=always

[Install]
WantedBy=multi-user.target